# system-config/cloudflared/cloudflared-furnace-vm.service
[Unit]
Description=Cloudflare Tunnel for VictoriaMetrics (furnace-vm)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/cloudflared tunnel --no-autoupdate --config /etc/cloudflared-furnace-vm/config.yml run
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only

[Install]
WantedBy=multi-user.target
