# system-config/cloudflared/cloudflared-furnace.service
[Unit]
Description=Cloudflare Tunnel for furnace display endpoint
After=network-online.target furnace-display.service
Wants=network-online.target
Requires=furnace-display.service

[Service]
Type=simple
ExecStart=/usr/local/bin/cloudflared tunnel --no-autoupdate --config /etc/cloudflared-furnace/config.yml run
Restart=always
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only

[Install]
WantedBy=multi-user.target
