# quadlets/influxdb.container
[Unit]
Description=InfluxDB v3 Container
After=network-online.target
Wants=network-online.target

[Container]
ContainerName=influxdb3
Image=docker.io/library/influxdb:3.3-core
PublishPort=0.0.0.0:8181:8181
Volume=%h/.influxdb3:/var/lib/influxdb3
Volume=%h/.influxdb3-plugins:/plugins
User=1000:1000
Exec=influxdb3 serve --object-store file --data-dir /var/lib/influxdb3 --node-id piiot

[Install]
WantedBy=default.target
