# quadlets/telegraf.container
[Unit]
Description=Telegraf Data Collection Agent
After=influxdb.service hivemq.service
Requires=influxdb.service hivemq.service

[Container]
ContainerName=telegraf
Image=docker.io/library/telegraf:latest
Network=host
Volume=%h/sensor-pipeline/configs/telegraf.conf:/etc/telegraf/telegraf.conf:ro
PodmanArgs=--privileged

[Service]
ExecStartPre=/bin/bash -c 'until nc -z 127.0.0.1 1883; do sleep 2; done'
ExecStartPre=/bin/bash -c 'modprobe bme680_i2c || true'
ExecStartPre=/bin/bash -c 'echo "bme680 0x77" > /sys/bus/i2c/devices/i2c-1/new_device || true'
ExecStartPre=/bin/bash -c 'sleep 1 && mkdir -p /dev/iio && ln -sf ../iio:device0 /dev/iio/bme680 || true'

[Install]
WantedBy=default.target
