# system-config/furnace-display/furnace-display.service
[Unit]
Description=Furnace display JSON endpoint
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=louthenw
Group=louthenw
EnvironmentFile=/etc/furnace-display.env
ExecStart=/usr/bin/python3 /opt/furnace-display/furnace_endpoint.py
Restart=always
RestartSec=3
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
StateDirectory=furnace-display

[Install]
WantedBy=multi-user.target
